# docker/hive2-no-auth/files/supervisord.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
startsecs=2
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/var/log/sshd.log
autostart=false
